# deploy/leapex.service
[Unit]
Description=Leapex FastAPI backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/leapex
Environment=PATH=/opt/leapex/.venv/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/leapex/.venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
